# systemd/set-theme@light.timer
[Unit]
Description=Run set-theme@light.service daily at 8:00

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
